# rc.local
#!/bin/sh -e
# Debian: /etc/rc.local

##  Protect sh and bash init scripts         ##
##  to prevent privilege escalation attacks  ##
chfiles="/home/user/.bashrc /home/user/.bash_profile /home/user/.bash_login /home/user/.bash_logout /home/user/.profile"
touch $chfiles || true
chown -f root:root $chfiles || true
chattr -f +i $chfiles || true


# end of script
exit 0
